# power-management/templates/wake-event.service.j2
# {{ ansible_managed }}
# VMStation Wake Event Handler Service
# Deployed to: /etc/systemd/system/vmstation-wake-event.service

[Unit]
Description=VMStation Wake Event Handler
Documentation=https://github.com/jjbly-vmstation/cluster-setup
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 {{ handler_dir }}/wake-handler.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# User and group (run as root for WoL access)
User=root
Group=root

# Environment
Environment="PYTHONUNBUFFERED=1"
{% if wake_auth_token is defined and wake_auth_token %}
Environment="VMSTATION_WAKE_TOKEN={{ wake_auth_token }}"
{% endif %}

# Working directory
WorkingDirectory={{ handler_dir }}

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ log_dir }}
PrivateTmp=yes

# Network access required
PrivateNetwork=no

# Resource limits
MemoryMax=256M
CPUQuota=20%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vmstation-wake-handler

[Install]
WantedBy=multi-user.target
