# deprecated/systemd/vmstation-wake-event.service
# VMStation Wake Event Handler Service
# This service provides an HTTP endpoint for wake-on-LAN event handling.
#
# Installation:
#   sudo cp vmstation-wake-event.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable vmstation-wake-event.service
#   sudo systemctl start vmstation-wake-event.service
#
# The service exposes endpoints:
#   GET  /health          - Health check
#   POST /wake/<hostname> - Wake specific node
#   POST /wake/all        - Wake all nodes
#   GET  /status          - Cluster power status

[Unit]
Description=VMStation Wake Event Handler
Documentation=https://github.com/jjbly-vmstation/cluster-setup
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/vmstation/wake-handler/wake-handler.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# User and group
# Run as root to allow sending WoL packets
User=root
Group=root

# Working directory
WorkingDirectory=/opt/vmstation/wake-handler

# Environment configuration
Environment="PYTHONUNBUFFERED=1"
# Set VMSTATION_WAKE_TOKEN for authentication
# Environment="VMSTATION_WAKE_TOKEN=your-secret-token"

# Security hardening (relaxed for network access)
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/vmstation
PrivateTmp=yes

# Network access required
PrivateNetwork=no

# Resource limits
MemoryMax=256M
CPUQuota=20%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vmstation-wake-handler

[Install]
WantedBy=multi-user.target
